# deploy/ansible/playbooks/roles/microservice/templates/supervisord/supervisord.conf
[program:{{app_name}}] 
user={{ user }}
command={{ app_dir }}/app/bin/uwsgi -i {{app_dir}}/conf/uwsgi/uwsgi.ini
directory={{ app_dir }}/app/src
autostart=true
autorestart=true
stopsignal=INT
